# deploy/ozon-dashboard.service
[Unit]
Description=Ozon Dashboard aiohttp app
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/ozon-dashboard
EnvironmentFile=/etc/ozon-dashboard.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONUTF8=1
ExecStart=/opt/ozon-dashboard/.venv/bin/python run_dashboard_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
