# deploy/ozon-price-control-sync.timer
[Unit]
Description=Run Ozon Seller price-control customer price sync twice daily

[Timer]
OnCalendar=*-*-* 08,20:00:00
Persistent=true
RandomizedDelaySec=300
Unit=ozon-price-control-sync.service

[Install]
WantedBy=timers.target
