# deploy/ozon-price-control-sync.timer
[Unit]
Description=Run Ozon Seller price-control customer price sync every hour

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300
Unit=ozon-price-control-sync.service

[Install]
WantedBy=timers.target
